# proj_1/resources/my_test_systemd.service
[Unit]
Description=My C++ Application
After=network.target

[Service]
ExecStart=/opt/test/test_systemd/test_systemd
Restart=always
StandardOutput=journal
StandardError=journal

WatchdogSec=30s
Restart=on-failure

[Install]
WantedBy=multi-user.target